# .config/systemd/user/xdg-desktop-portal-wlr.service.d/override.conf
[Service]
Environment=XDG_CURRENT_DESKTOP=sway
Environment=WAYLAND_DISPLAY=wayland-1
