# .local/bin/postinstall/root/etc/systemd/system/audio-switch.service
[Unit]
Description=Set default ALSA audio device
After=sound.target

[Service]
Type=oneshot
User=pavel
ExecStart=/home/pavel/.local/bin/audio-switch -d
KillMode=process

[Install]
WantedBy=multi-user.target
