# .local/bin/postinstall/root/etc/systemd/system/systemd-suspend.service.d/chronyd.conf
[Service]
ExecStartPost=/usr/bin/chronyd -q
